# Resources/Rasp4.timer
[Unit]
Description=4040 Service Timer

[Timer]
OnBootSec=1min
Unit=Rasp4.Service

[Install]
WantedBy=default.target
