# Resources/Raspi_Daemon.timer
[Unit]
Description=1010 Service Timer

[Timer]
OnBootSec=1min
Unit=1010.service

[Install]
WantedBy=default.target